# auralis-enhanced.service
[Unit]
Description=Auralis Enhanced OpenAI-compatible TTS API
After=network-online.target nvidia-persistenced.service
Wants=network-online.target

[Service]
Type=simple
User=op
WorkingDirectory=/home/op/Auralis-Enhanced
Environment=CUDA_DEVICE_ORDER=PCI_BUS_ID
# This host has three CUDA GPUs (0-2); use the otherwise-free third GPU by default.
# /etc/default/auralis-enhanced may override this and any other Environment value.
Environment=CUDA_VISIBLE_DEVICES=2
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/auralis-enhanced
ExecStartPre=/usr/bin/test -x /home/op/miniconda3/envs/auralis-enhanced/bin/auralis.openai
ExecStart=/home/op/miniconda3/envs/auralis-enhanced/bin/auralis.openai --host 0.0.0.0 --port 6688 --backend vllm --device cuda --gpu_memory_utilization 0.65 --cpu_offload_gb 0 --swap_space 0 --inactivity_timeout 300 --cleanup_interval 30
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
